# port_monit.service
[Unit]
Description=Port Monitor Service

[Service]
ExecStart=/bin/bash ~/port_monit.sh
